# os-sim/contrib/mrtg/mrtg-ossim.cron
0-59/5 * * * * /opt/ossim/mrtg/launch-mrtg

