# ossim-cd-tools/overrides/munin-node
#!/bin/bash
### BEGIN INIT INFO
# Provides:          munin-node
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: munin-node
# Description:
#
### END INIT INFO

