# alienvault-dummies-config/debian/alienvault-ossec.cron.daily
#!/bin/bash
find /var/ossec/logs/alerts/ -mmin +1440 -name "*.log.*" -type f -exec rm {} \;
cd /var/ossec/logs/alerts/ && find . -type d -empty -delete

chown www-data:ossec /var/ossec/logs/ -R
chmod ug+w /var/ossec/logs -R
