# os-sim/etc/cron.daily/check-tickets
#!/bin/bash
/usr/bin/php /usr/share/ossim/scripts/incidents/check_tickets.php
