# os-sim/etc/cron.daily/clean-alarms
#!/bin/bash
# DELETE ALARMS
/usr/bin/php /usr/share/ossim/scripts/alarms/delete_alarms.php
# ASSET SYNC DUMP
/usr/share/ossim/scripts/assets_sync.sh
